# files/profile.d/docker_profile.sh
#!/bin/bash

############
# Command to get the id of the last used docker container
# ----------
alias dlast='docker ps -l -q'

